# container/files/etc/supervisor/conf.d/supervisor-app.conf
[program:app]
command = /var/www/venv/bin/gunicorn -c /etc/app/conf.py myapp:app
autostart = true
autorestart = true
startsec = 5
redirect_stderr = true
