# container/files/etc/supervisor/conf.d/supervisor-nginx.conf
[program:nginx]
command = /usr/sbin/nginx
autostart = true
autorestart = true
startsec = 10
redirect_stderr = true
